# services/vpnclient.service
[Unit]
Description=SoftEther VPN Client
After=network.target

[Service]
ExecStart=/usr/local/vpnclient/vpnclient start
ExecStop=/usr/local/vpnclient/vpnclient stop
Restart=on-abort
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target